# LAB8/watchlog.timer
[Unit]
Description=Run Watchlog Script every 10 seconds

[Timer]
OnCalendar=*:*:0/10
Unit=watchlog.service

[Install]
WantedBy=multi-user.target
