# LAB8/watchlog.timer
[Unit]
Description=Run Watchlog Script every 5 seconds

[Timer]
OnBootSec=5
Unit=watchlog.service

[Install]
WantedBy=multi-user.target
