# LAB8/watchlog.service
[Unit]
Description=Watchlog Service

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/watchlog
ExecStart=/root/watchlog.sh $WORD $LOG
