# LAB8/watchlog.timer
[Unit]
Description=Run Watchlog Script every 30 seconds

[Timer]
OnBootSec=30
Unit=watchlog.service

[Install]
WantedBy=multi-user.target
